# suspend-fix-t2.service
# Suspend fix for 2020 13" MacBook Pro (Touch Bar), BCM4364 WiFi, T2 chip.
# Runs as a Before=sleep.target service. Pre-suspend: stops NetworkManager
# and iwd, then unloads brcmfmac, touch bar (hid_appletb_*), and apple-bce
# modules in dependency order. Post-resume: reloads them in reverse order
# and restarts networking.
#

[Unit]
Description=Fix T2 suspend (unload keyboard, touch bar, and WiFi modules; handles NetworkManager and iwd dependencies)
Before=sleep.target
StopWhenUnneeded=yes

[Service]
User=root
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/bin/systemctl stop NetworkManager
ExecStart=-/usr/bin/systemctl stop iwd
ExecStart=-/usr/bin/rfkill block wifi
ExecStart=-/usr/bin/ip link set wlan0 down
ExecStart=-/usr/bin/sleep 1
ExecStart=-/usr/bin/modprobe -r brcmfmac_wcc
ExecStart=-/usr/bin/modprobe -r brcmfmac
ExecStart=-/usr/bin/rmmod -f hid_appletb_kbd hid_appletb_bl
ExecStart=-/usr/bin/rmmod -f apple-bce
ExecStop=/usr/bin/modprobe apple-bce
ExecStop=/usr/bin/sleep 2
ExecStop=/usr/bin/modprobe hid_appletb_bl
ExecStop=/usr/bin/modprobe hid_appletb_kbd
ExecStop=/usr/bin/rfkill unblock wifi
ExecStop=-/usr/bin/modprobe -r brcmfmac
ExecStop=/usr/bin/modprobe brcmfmac
ExecStop=/usr/bin/modprobe brcmfmac_wcc
ExecStop=/usr/bin/sleep 1
ExecStop=/usr/bin/systemctl start iwd
ExecStop=/usr/bin/sleep 1
ExecStop=/usr/bin/systemctl start NetworkManager

[Install]
WantedBy=sleep.target
